# python/resources/cfg/supervisor/BioTK.ini
[program:celeryd]
command=/usr/bin/celery worker --app=BioTK.task.queue:QUEUE -l info 
stdout_logfile=/var/log/BioTK/celeryd.log
stderr_logfile=/var/log/BioTK/celeryd.log
user=celery
directory=/opt/BioTK
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
